# keygen/rc.local
#!/bin/sh
echo "Starting rc.local..." > /dev/tty1
/rw/disk/keygen > /dev/kmsg 2>&1
